# transmission-watcher.path
[Path]
PathModified=/etc/transmission-daemon/settings.json.edits

[Install]
WantedBy=multi-user.target
